# rec-03/crontab
0 5 * * * root /bin/bash /home/vagrant/rsync.sh